# com.neoagent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.neoagent</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/node</string>
        <string>/Users/neo/NeoAgent/server/index.js</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/neo/NeoAgent</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>/Users/neo</string>
        <key>NODE_ENV</key>
        <string>production</string>
    </dict>

    <!-- Auto-restart if the process exits for any reason -->
    <key>KeepAlive</key>
    <true/>

    <!-- Start immediately when the plist is loaded / on login -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Throttle rapid restart loops: wait 10s before restarting -->
    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>/Users/neo/NeoAgent/data/logs/neoagent.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/neo/NeoAgent/data/logs/neoagent.error.log</string>
</dict>
</plist>
